# deploy/systemd/mainsbot-web@.service.example
[Unit]
Description=MainsBot web instance %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/MainsBot
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /opt/MainsBot/main-web.js --config /opt/MainsBot/config/%i.ini
Restart=always
RestartSec=3

# Per-instance runtime dir for unix sockets:
RuntimeDirectory=mainsbot/%i
RuntimeDirectoryMode=0775

[Install]
WantedBy=multi-user.target
