# deploy/systemd/mainsbot@.service.example
[Unit]
Description=MainsBot instance %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/MainsBot
Environment=NODE_ENV=production
Environment=MODULE_WEB=0
ExecStart=/usr/bin/node /opt/MainsBot/main.js --config /opt/MainsBot/config/%i.ini
Restart=always
RestartSec=3

# Bot-only unit. Do not own the shared /run/mainsbot socket directory here,
# otherwise stopping one bot instance can unlink another instance's web socket.

[Install]
WantedBy=multi-user.target
